# softicar-github-runner.service
# A systemd service that keeps the ephemeral SoftiCAR GitHub Runner alive.
# 
# Install with:
#   sudo install -m 644 softicar-github-runner.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable softicar-github-runner
# Start with:
#   sudo systemctl start softicar-github-runner
# Stop with:
#   sudo systemctl stop softicar-github-runner
# See live logs with:
#   journalctl -f -u softicar-github-runner --no-hostname --no-tail

[Unit]
Description=SoftiCAR GitHub Runner Container
After=docker.service
Requires=docker.service

[Service]
TimeoutStartSec=0
Restart=always
User=softicar
ExecStart=/home/softicar/softicar-github-runner/softicar-github-runner.sh

[Install]
WantedBy=multi-user.target
